# libComposite.service
[Unit]
Description=Start libComposite

[Service]
Type=oneshot
ExecStart=/root/usb.sh

[Install]
WantedBy=multi-user.target
